# stow/xdg_config/zsh/profile.d/cargo.sh
# cargo
if [ -d "$XDG_DATA_HOME"/cargo ]; then
    export CARGO_HOME="$XDG_DATA_HOME"/cargo
    source $CARGO_HOME/env
fi
